# roles/bgp/templates/birdnotify.service.j2
[Unit]
Description=React to a path unit for Bird VIP advertisement

[Service]
Type=oneshot
ExecStart={{ birdwatch_script_path }}
